# packaging/deb/sessionguard-agent.service
[Unit]
Description=SessionGuard Linux Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sessionguard/agent
Environment=DOTNET_ENVIRONMENT=Production
EnvironmentFile=-/etc/sessionguard/agent.env
ExecStart=/usr/bin/dotnet /opt/sessionguard/agent/Agent.Linux.dll
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
