# docs/sessionguard-agent.service
[Unit]
Description=SessionGuard Linux Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sessionguard
ExecStart=/usr/bin/dotnet /opt/sessionguard/Agent.Linux.dll
Restart=always
RestartSec=10
Environment=DOTNET_ENVIRONMENT=Production
# Root is used in the MVP so the agent can resolve and lock desktop sessions via loginctl.

[Install]
WantedBy=multi-user.target
